# framework/scripts/agent-activity-tracker.service
[Unit]
Description=TTA.dev Agent Activity Tracker
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=thein
WorkingDirectory=/home/thein/repos/TTA.dev
ExecStart=/home/thein/.local/bin/uv run python scripts/agent-activity-tracker-tta.py --workspace /home/thein/repos/TTA.dev --port 8001
Restart=always
RestartSec=10
StandardOutput=append:/var/log/tta-agent-tracker.log
StandardError=append:/var/log/tta-agent-tracker.log

# Environment
Environment="PATH=/home/thein/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/thein/repos/TTA.dev/packages"

[Install]
WantedBy=multi-user.target
